# schedule.cron
58 * * * * python3 /usr/local/oci-scheduler/AutoScaleALL.py -a Down
01 * * * * python3 /usr/local/oci-scheduler/AutoScaleALL.py -a Up
