# schedule.cron
58 * * * * python3 /usr/local/oci-scheduler/AutoScaleALL.py -ip -a Down
01 * * * * python3 /usr/local/oci-scheduler/AutoScaleALL.py -ip -a Up
